# src/root/etc/init.d/ntpd
#!/sbin/openrc-run

description="NTP daemon"
pidfile=/var/run/${RC_SVCNAME}.pid
supervisor=supervise-daemon

command="/sbin/ntpd"
command_args_foreground="-n"

depend() {
    need net-online
}
